# systemd/librespot.service
[Unit]
Description=librespot (Spotify Connect receiver for mini-vinyl)
After=bluetooth.target pipewire.service pipewire-pulse.service

[Service]
Type=simple
# --name must match SPOTIFY_DEVICE_NAME in config/secrets.env
# --backend pipewire routes audio straight to PipeWire, which owns the
# Bluetooth A2DP sink on Raspberry Pi OS Bookworm/Trixie. If your librespot
# build doesn't have the pipewire backend compiled in, use
# `--backend alsa --device default` instead (PipeWire intercepts ALSA's
# "default" device automatically once pipewire-alsa is installed).
ExecStart=/usr/bin/librespot \
  --name mini-vinyl \
  --backend pipewire \
  --bitrate 160 \
  --initial-volume 80
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
