# systemd/mini-vinyl.service
[Unit]
Description=mini-vinyl NFC record player
After=bluetooth.target pipewire.service network-online.target

[Service]
Type=simple
WorkingDirectory=%h/Projects/mini-vinyl
ExecStart=%h/Projects/mini-vinyl/venv/bin/python -m mini_vinyl.main
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
